# src/etc/systemd/system/kubeconfig.path
[Unit]
Description=Waits for the kubeconf generated by k0s

[Path]
PathExists=/var/lib/k0s/pki/admin.conf
Unit=kubeconfig.service

[Install]
WantedBy=multi-user.target
