# src/etc/systemd/system/kubeconfig.path
[Unit]
Description=Sets permissions on the kubeconf generated by k0s

[Path]
PathChanged=/var/lib/k0s/pki/admin.conf
Unit=kubeconfig.service

[Install]
WantedBy=multi-user.target
