# src/etc/systemd/system/kubeconfig.service
[Unit]
Description=Sets permissions on the kubeconf generated by k0s

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=chown :wheel /var/lib/k0s/pki/admin.conf
ExecStart=chmod g+r /var/lib/k0s/pki/admin.conf

[Install]
WantedBy=multi-user.target
